# etc/systemd/system/mnt-NAS.automount
[Unit]
Description=Automount NAS

[Automount]
Where=/mnt/NAS

[Install]
WantedBy=multi-user.target
